# luci-app-easymesh/root/etc/hotplug.d/iface/95-easymesh-healer
#!/bin/sh
[ "$INTERFACE" = "lan" ] || [ "$INTERFACE" = "wan" ] || exit 0
ip link show bat0 >/dev/null 2>&1 || exit 0
[ "$(uci -q get easymesh.global.enabled)" = "1" ] || exit 0
MESH_IF="bat0"
BATMAN_OGM_NORMAL=500
EM_RUN_DIR="/tmp/easymesh"
LOGFILE="/tmp/easymesh.log"

_emlog() {
    local ts; ts=$(date "+%Y-%m-%d %H:%M:%S")
    printf "%s [%-5s] [%-8s] %s\n" "$ts" "$1" "$2" "$3" >> "$LOGFILE"
}

case "$ACTION" in
    ifdown)
        batctl meshif "$MESH_IF" orig_interval 300 2>/dev/null
        _emlog "WARN" "heal" "LAN link down — OGM interval set to 300ms for fast wireless failover"
        logger -t EasyMesh "LAN link down, accelerating wireless failover" 2>/dev/null
        ;;
    ifup)
        batctl meshif "$MESH_IF" orig_interval "$BATMAN_OGM_NORMAL" 2>/dev/null
        _emlog "INFO" "heal" "LAN link up — OGM interval restored to ${BATMAN_OGM_NORMAL}ms"
        logger -t EasyMesh "LAN link restored, OGM interval back to normal" 2>/dev/null
        ;;
esac
