# luci-app-diskman-js/root/etc/uci-defaults/luci-app-diskman-js
#!/bin/sh
chmod +x /usr/libexec/diskman

rm -rf /tmp/luci-*
exit 0
