# luci-app-mesh-node/root/etc/uci-defaults/80_mesh-node
#!/bin/sh
echo "config mesh_node 'main'" > /etc/config/mesh_node
chmod +x /etc/init.d/mesh-node
exit 0
