# luci-app-mesh-node/root/etc/uci-defaults/luci-app-mesh-node
#!/bin/sh

echo "config mesh_node 'main'" > /etc/config/mesh_node
chmod +x /etc/init.d/mesh_node
rm -f /tmp/luci-*
exit 0
