# luci-app-mesh11sd/root/etc/uci-defaults/luci-app-mesh11sd
#!/bin/sh

/etc/init.d/mesh11sd stop    2>/dev/null || true
/etc/init.d/mesh11sd disable 2>/dev/null || true
rm -f /tmp/luci-*
exit 0
