# luci-app-easymesh/root/etc/init.d/easymesh
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=95
STOP=10

. /lib/functions/network.sh

UCI_APP=easymesh
LOGFILE="/tmp/easymesh.log"
lang=$(uci get luci.main.lang 2>/dev/null || echo "en")

translate() {
    local tpl="$1"; shift
    local translated
    translated=$(lua -e "require \"luci.i18n\".setlanguage(\"${lang}\"); print(require \"luci.i18n\".translate([==[${tpl}]==]))" 2>/dev/null)
    [ -z "$translated" ] && translated="$tpl"
    if [ $# -gt 0 ]; then
        printf "$translated" "$@"
    else
        printf '%s' "$translated"
    fi
}

_emlog() {
    local level="$1" tag="$2"; shift 2
    local ts; ts=$(date "+%Y-%m-%d %H:%M:%S")
    printf "%s [%-5s] [%-8s] %s\n" "$ts" "$level" "$tag" "$*" >> "$LOGFILE"
}
log()  { _emlog "INFO"  "initd" "$(translate "$@")"; }
loge() { _emlog "ERROR" "initd" "$(translate "$@")"; }

# ── 获取指定频段的 radio 设备名 ──────────────────────────────────────────────
_band_target=""; _found=""
_find_radio() {
    local _n="$1" _b
    config_get _b "$_n" band
    case "$_b" in
        2g|2.4GHz|bgn|bg|b) [ "$_band_target" = "2g" ] && _found="$_n" ;;
        5g|5GHz|an|a)        [ "$_band_target" = "5g" ] && _found="$_n" ;;
    esac
}
get_radio() {
    local _v="$1" _b="$2"
    _found=""; _band_target="$_b"
    config_load wireless
    config_foreach _find_radio wifi-device
    eval "$_v='$_found'"
}

_del_managed() {
    local _n="$1" _m
    config_get _m "$_n" easymesh_managed
    [ "$_m" = "1" ] && uci -q del wireless."$_n"
}
cleanup_wireless() {
    config_load wireless
    config_foreach _del_managed wifi-iface
    uci commit wireless
}

_agent_prepare() {
    local my_ip prefix suffix new_ip

    my_ip=$(ip -4 addr show dev br-lan 2>/dev/null \
        | awk '/inet /{print $2}' | cut -d/ -f1 | head -1)
    [ -z "$my_ip" ] && { log "agent_prepare: br-lan has no IP"; return 1; }

    suffix=$(echo "$my_ip" | cut -d. -f4)
    if [ "$suffix" != "1" ]; then
        log "agent_prepare: IP=%s (not .1, no change needed)" "$my_ip"
        return 0
    fi

    prefix=$(echo "$my_ip" | cut -d. -f1-3)
    log "agent_prepare: IP=%s (.1 conflict detected), computing stable IP..." "$my_ip"

    # 用 br-lan MAC 后3字节 → 十六进制 → 十进制 → 模98 → 加100基准
    local mac_raw
    mac_raw=$(ip link show br-lan 2>/dev/null | awk '/ether/{print $2}' | tr -d ':')
    [ -z "$mac_raw" ] && mac_raw=$(cat /sys/class/net/eth0/address 2>/dev/null | tr -d ':')
    [ -z "$mac_raw" ] && mac_raw="000000aabbcc"
    local mac_tail; mac_tail=$(printf '%s' "$mac_raw" | tail -c 6)
    local mac_dec;  mac_dec=$(printf '%d' "0x${mac_tail}" 2>/dev/null || echo 0)
    local offset;   offset=$(( mac_dec % 98 ))
    local candidate="${prefix}.$((100 + offset))"

    new_ip="$candidate"
    # 碰撞检测：如果派生IP已被占用则线性扫描
    if ping -c 1 -W 1 -q "$candidate" >/dev/null 2>&1; then
        log "agent_prepare: candidate %s busy, linear scan..." "$candidate"
        new_ip=""
        local i
        for i in $(seq 100 197); do
            local c="${prefix}.${i}"
            ping -c 1 -W 1 -q "$c" >/dev/null 2>&1 || { new_ip="$c"; break; }
        done
        [ -z "$new_ip" ] && {
            loge "agent_prepare: no free IP in .100-.197, keeping .1"
            return 0
        }
    fi

    log "agent_prepare: assigning %s -> %s" "$my_ip" "$new_ip"
    uci set network.lan.ipaddr="$new_ip"
    uci commit network
    uci set system.@system[0].hostname="OpenWrt-Agent-$(echo "$new_ip" | cut -d. -f4)"
    uci commit system
    /etc/init.d/network reload

    local _w=0
    while [ $_w -lt 20 ]; do
        my_ip=$(ip -4 addr show dev br-lan 2>/dev/null \
            | awk '/inet /{print $2}' | cut -d/ -f1 | head -1)
        [ "$my_ip" = "$new_ip" ] && break
        sleep 1; _w=$((_w+1))
    done
    log "agent_prepare: done. br-lan IP=%s" "$my_ip"
}

_classify_radio() {
    local _n="$1" _b
    config_get _b "$_n" band
    case "$_b" in
        2g*|2.4*|bgn|bg|b)
            [ -z "$radio_2g" ] && radio_2g="$_n" ;;
        5g*|5GHz*|an|a)
            _seen5g=$((_seen5g+1))
            if [ "$_seen5g" -eq 1 ]; then
                [ -z "$radio_5g" ] && radio_5g="$_n"
            else
                [ -z "$radio_5g2" ] && radio_5g2="$_n"
            fi
            ;;
    esac
}

_setup_ap_iface() {
    local name="$1" radio="$2" ssid="$3" key="$4"
    local ieee80211r="$5" mobility_domain="$6" ft_over_ds="$7"
    local reassociation_deadline="$8" ieee80211k="$9" ieee80211v="${10}"
    uci set wireless.${name}=wifi-iface
    uci set wireless.${name}.device="$radio"
    uci set wireless.${name}.mode='ap'
    uci set wireless.${name}.ssid="$ssid"
    uci set wireless.${name}.network='lan'
    uci set wireless.${name}.easymesh_managed='1'
    uci set wireless.${name}.disabled='0'
    if [ -n "$key" ]; then
        uci set wireless.${name}.encryption='psk2+ccmp'
        uci set wireless.${name}.key="$key"
    else
        uci set wireless.${name}.encryption='none'
    fi
    [ "$ieee80211r" = "1" ] && {
        uci set wireless.${name}.ieee80211r='1'
        uci set wireless.${name}.mobility_domain="$mobility_domain"
        uci set wireless.${name}.ft_over_ds="$ft_over_ds"
        uci set wireless.${name}.ft_psk_generate_local='1'
        uci set wireless.${name}.reassociation_deadline="$reassociation_deadline"
    }
    [ "$ieee80211k" = "1" ] && {
        uci set wireless.${name}.ieee80211k='1'
        uci set wireless.${name}.rrm_neighbor_report='1'
        uci set wireless.${name}.rrm_beacon_report='1'
    }
    [ "$ieee80211v" = "1" ] && {
        uci set wireless.${name}.ieee80211v='1'
        uci set wireless.${name}.bss_transition='1'
    }
}

apply_mesh_config() {
    local role backhaul mesh_id mesh_key mesh_band ssid key
    local ieee80211r mobility_domain ft_over_ds reassociation_deadline
    local ieee80211k ieee80211v routing_algo mesh_fwding

    config_load "$UCI_APP"
    config_get role                   global role                  'master'
    config_get backhaul               global backhaul              'wired'
    config_get mesh_id                global mesh_id               'OpenWrt-Mesh'
    config_get mesh_key               global mesh_key              ''
    config_get mesh_band              global mesh_band             '5g'
    config_get ssid                   global ssid                  'OpenWrt'
    config_get key                    global key                   ''
    config_get ieee80211r             global ieee80211r            '1'
    config_get mobility_domain        global mobility_domain       'aabb'
    config_get ft_over_ds             global ft_over_ds            '1'
    config_get reassociation_deadline global reassociation_deadline '1000'
    config_get ieee80211k             global ieee80211k            '1'
    config_get ieee80211v             global ieee80211v            '1'
    config_get routing_algo           global routing_algo          'BATMAN_IV'
    config_get mesh_fwding            global mesh_fwding           '0'

    cleanup_wireless

    local radio_2g="" radio_5g="" radio_5g2=""
    local _seen5g=0
    config_load wireless
    config_foreach _classify_radio wifi-device
    local radio_mesh
    if [ "$mesh_band" = "2g" ]; then
        radio_mesh="$radio_2g"
    else
        radio_mesh="$radio_5g"
    fi
    [ -z "$radio_mesh" ] && radio_mesh=$(uci show wireless 2>/dev/null | awk -F= '/=wifi-device/{print $1}' | head -1 | sed 's/wireless\.//')
    log "apply_mesh_config: role=%s backhaul=%s mesh_radio=%s 2g=%s 5g=%s" "$role" "$backhaul" "$radio_mesh" "$radio_2g" "$radio_5g"

    # ── network.mesh dummy section ────────────────────────────────────────────
    uci set network.mesh=interface
    uci set network.mesh.proto='none'
    uci commit network

    # ── 无线 mesh 回程（802.11s）─────────────────────────────────────────────
    uci set wireless.mesh0=wifi-iface
    uci set wireless.mesh0.device="$radio_mesh"
    uci set wireless.mesh0.mode='mesh'
    uci set wireless.mesh0.mesh_id="$mesh_id"
    uci set wireless.mesh0.mesh_fwding="$mesh_fwding"
    uci set wireless.mesh0.mesh_rssi_threshold='-80'
    uci set wireless.mesh0.network='mesh'
    uci set wireless.mesh0.easymesh_managed='1'
    if [ -n "$mesh_key" ]; then
        uci set wireless.mesh0.encryption='sae'
        uci set wireless.mesh0.key="$mesh_key"
    else
        uci set wireless.mesh0.encryption='none'
    fi
    uci set wireless.mesh0.disabled=$([ "$backhaul" = "wired" ] && echo '1' || echo '0')

    # ── DHCP server 策略 ──────────────────────────────────────────────────────
    if [ "$role" = "master" ]; then
        uci -q del dhcp.lan.ignore
        log "Role: master — DHCP enabled"
    else
        uci set dhcp.lan.ignore='1'
        log "Role: agent — DHCP disabled (master handles all address allocation)"
    fi
    uci commit dhcp

    # ── 客户端 AP（双频合并）───────────────────────────────
    local _ap_idx=0
    # 5G AP（主）
    if [ -n "$radio_5g" ]; then
        _setup_ap_iface easymesh_ap "$radio_5g" "$ssid" "$key" "$ieee80211r" "$mobility_domain" "$ft_over_ds" "$reassociation_deadline" "$ieee80211k" "$ieee80211v"
        _ap_idx=$((_ap_idx+1))
        log "AP created: easymesh_ap on %s (5G, ssid=%s)" "$radio_5g" "$ssid"
    fi
    # 2.4G AP（补充覆盖）
    if [ -n "$radio_2g" ]; then
        _setup_ap_iface easymesh_ap_2g "$radio_2g" "$ssid" "$key" "$ieee80211r" "$mobility_domain" "$ft_over_ds" "$reassociation_deadline" "$ieee80211k" "$ieee80211v"
        _ap_idx=$((_ap_idx+1))
        log "AP created: easymesh_ap_2g on %s (2.4G, ssid=%s)" "$radio_2g" "$ssid"
    fi
    # 兜底：如果没有识别到任何 radio，用 radio0
    if [ "$_ap_idx" -eq 0 ]; then
        _setup_ap_iface easymesh_ap "radio0" "$ssid" "$key" "$ieee80211r" "$mobility_domain" "$ft_over_ds" "$reassociation_deadline" "$ieee80211k" "$ieee80211v"
        log "AP created: easymesh_ap on radio0 (fallback)"
    fi
    # 确保相关 radio 设备本身不处于 disabled 状态
    [ -n "$radio_mesh" ] && uci -q del wireless.${radio_mesh}.disabled
    [ -n "$radio_5g" ]   && uci -q del wireless.${radio_5g}.disabled
    [ -n "$radio_2g" ]   && uci -q del wireless.${radio_2g}.disabled
    uci commit wireless
    log "apply_mesh_config done (%s AP(s) created)" "$_ap_idx"
}

# ── procd ─────────────────────────────────────────────────────────────────────
_open_firewall_ports() {
    if command -v iptables >/dev/null 2>&1; then
        iptables -C INPUT -p tcp --dport 4304 -j ACCEPT 2>/dev/null \
            || iptables -I INPUT 1 -p tcp --dport 4304 -j ACCEPT 2>/dev/null
        iptables -C INPUT -p udp --dport 4305 -j ACCEPT 2>/dev/null \
            || iptables -I INPUT 1 -p udp --dport 4305 -j ACCEPT 2>/dev/null
        log "Firewall: opened TCP/UDP 4304/4305 via iptables"
    fi
    if command -v nft >/dev/null 2>&1; then
        nft insert rule inet fw4 input tcp dport 4304 accept 2>/dev/null \
            || nft insert rule inet filter input tcp dport 4304 accept 2>/dev/null || true
        nft insert rule inet fw4 input udp dport 4305 accept 2>/dev/null \
            || nft insert rule inet filter input udp dport 4305 accept 2>/dev/null || true
        log "Firewall: opened TCP/UDP 4304/4305 via nft"
    fi
}

start_service() {
    local enabled role backhaul
    config_load "$UCI_APP"
    config_get_bool enabled global enabled 0
    [ "$enabled" -eq 0 ] && return 0

    config_get role     global role     'master'
    config_get backhaul global backhaul 'wired'
    log "EasyMesh starting (role=%s)" "$role"

    [ "$role" = "agent" -a ! -f /etc/easymesh_agent ] && {
        local ipaddr; ipaddr=$(uci_get network lan ipaddr '192.168.1.1')
        local hostname; hostname=$(uci_get system @system[0] hostname openwrt)
        printf 'ipaddr=%s\nhostname=%s\n' "$ipaddr" "$hostname" > /etc/easymesh_agent
    }

    apply_mesh_config

    if [ "$role" = "agent" ] && [ "$backhaul" = "wired" ]; then
        _agent_prepare
    fi

    _open_firewall_ports

    local _lip _lw=0
    while [ $_lw -lt 20 ]; do
        _lip=$(ip -4 addr show dev br-lan 2>/dev/null | awk '/inet /{print $2}' | cut -d/ -f1 | head -1)
        [ -n "$_lip" ] && break
        sleep 1; _lw=$((_lw+1))
    done
    if [ -z "$_lip" ]; then
        log "WARN: br-lan has no IP after 20s, starting daemon anyway"
    else
        log "br-lan IP ready: %s" "$_lip"
    fi

    log "Reloading wireless config..."
    wifi reload
    log "Wireless reload done"

    log "Starting EasyMesh daemon..."
    procd_open_instance "easymesh-${role}"
    procd_set_param command /usr/sbin/easymesh-${role}
    procd_set_param respawn 10 3 0
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance
}

_close_firewall_ports() {
    if command -v iptables >/dev/null 2>&1; then
        while iptables -C INPUT -p tcp --dport 4304 -j ACCEPT 2>/dev/null; do
            iptables -D INPUT -p tcp --dport 4304 -j ACCEPT 2>/dev/null || break
        done
        while iptables -C INPUT -p udp --dport 4305 -j ACCEPT 2>/dev/null; do
            iptables -D INPUT -p udp --dport 4305 -j ACCEPT 2>/dev/null || break
        done
    fi
}

stop_service() {
    log "EasyMesh stopping..."

    for pid in $(pgrep -f "easymesh-(master|agent)" 2>/dev/null); do
        [ "$pid" != "$$" ] && kill "$pid" 2>/dev/null || true
    done

    # 停止无线 onboard AP（临时开放热点）
    if [ -f /tmp/easymesh_pair_hostapd.pid ]; then
        kill "$(cat /tmp/easymesh_pair_hostapd.pid)" 2>/dev/null || true
    fi
    if [ -f /tmp/easymesh_pair_dnsmasq.pid ]; then
        kill "$(cat /tmp/easymesh_pair_dnsmasq.pid)" 2>/dev/null || true
    fi
    iw dev easymesh-pair0 del 2>/dev/null || true

    _close_firewall_ports
    cleanup_wireless

    uci -q del dhcp.lan.ignore
    uci commit dhcp

    uci -q del network.mesh
    uci commit network
    ( sleep 1; wifi reload ) &
    rm -rf /tmp/easymesh* /tmp/em_* 2>/dev/null
}

reload_service() {
    stop
    sleep 2
    start
}

service_triggers() {
    procd_add_reload_trigger "$UCI_APP"
}
