# luci-app-easymesh/root/etc/init.d/easymesh
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=95
STOP=10

. /lib/functions/network.sh

UCI_APP=easymesh
LOGFILE="/tmp/easymesh.log"
lang=$(uci get luci.main.lang 2>/dev/null || echo "en")

translate() {
    local tpl="$1"; shift
    local translated
    translated=$(lua -e "require \"luci.i18n\".setlanguage(\"${lang}\"); print(require \"luci.i18n\".translate([==[${tpl}]==]))" 2>/dev/null)
    [ -z "$translated" ] && translated="$tpl"
    if [ $# -gt 0 ]; then
        printf "$translated" "$@"
    else
        printf '%s' "$translated"
    fi
}

_emlog() {
    local level="$1" tag="$2"; shift 2
    local ts; ts=$(date "+%Y-%m-%d %H:%M:%S")
    printf "%s [%-5s] [%-8s] %s\n" "$ts" "$level" "$tag" "$*" >> "$LOGFILE"
}
log()  { _emlog "INFO"  "initd" "$(translate "$@")"; }
loge() { _emlog "ERROR" "initd" "$(translate "$@")"; }

# ── 获取指定频段的 radio 设备名 ──────────────────────────────────────────────
_band_target=""; _found=""
_find_radio() {
    local _n="$1" _b
    config_get _b "$_n" band
    case "$_b" in
        2g|2.4GHz|bgn|bg|b) [ "$_band_target" = "2g" ] && _found="$_n" ;;
        5g|5GHz|an|a)        [ "$_band_target" = "5g" ] && _found="$_n" ;;
    esac
}
get_radio() {
    local _v="$1" _b="$2"
    _found=""; _band_target="$_b"
    config_load wireless
    config_foreach _find_radio wifi-device
    eval "$_v='$_found'"
}

_del_managed() {
    local _n="$1" _m
    config_get _m "$_n" easymesh_managed
    [ "$_m" = "1" ] && uci -q del wireless."$_n"
}
cleanup_wireless() {
    config_load wireless
    config_foreach _del_managed wifi-iface
    uci commit wireless
}

_agent_prepare() {
    local my_ip prefix new_ip i

    my_ip=$(ip -4 addr show dev br-lan 2>/dev/null \
        | awk '/inet /{print $2}' | cut -d/ -f1 | head -1)
    [ -z "$my_ip" ] && { log "agent_prepare: br-lan has no IP"; return 1; }
    prefix=$(echo "$my_ip" | cut -d. -f1-3)

    if [ "$(echo "$my_ip" | cut -d. -f4)" = "1" ]; then
        log "agent_prepare: IP=%s (.1), finding free IP..." "$my_ip"
        new_ip=""
        for i in $(seq 100 254); do
            local c="${prefix}.${i}"
            ping -c 1 -W 1 -q "$c" >/dev/null 2>&1 || { new_ip="$c"; break; }
        done
        if [ -z "$new_ip" ]; then
            loge "agent_prepare: no free IP in .100-.254, keeping .1"
        else
            log "agent_prepare: assigning %s -> %s" "$my_ip" "$new_ip"
            uci set network.lan.ipaddr="$new_ip"
            uci commit network
            uci set system.@system[0].hostname="OpenWrt-Agent-$(echo "$new_ip" | cut -d. -f4)"
            uci commit system
            /etc/init.d/network reload
            local _w=0
            while [ $_w -lt 15 ]; do
                my_ip=$(ip -4 addr show dev br-lan 2>/dev/null \
                    | awk '/inet /{print $2}' | cut -d/ -f1 | head -1)
                [ "$my_ip" = "$new_ip" ] && break
                sleep 1; _w=$((_w+1))
            done
            log "agent_prepare: br-lan IP now %s" "$my_ip"
        fi
    else
        log "agent_prepare: IP=%s (not .1), skipping IP change" "$my_ip"
    fi

    # ── 步骤2：设置唯一MAC ────────────────
    local saved_mac
    saved_mac=$(uci -q get easymesh.global.agent_br_mac 2>/dev/null)

    if [ -z "$saved_mac" ]; then
        local seed
        seed=$(cat /sys/bus/soc/devices/soc0/serial_number 2>/dev/null)
        [ -z "$seed" ] && seed=$(grep -i "serial" /proc/cpuinfo 2>/dev/null | awk '{print $NF}' | head -1)
        [ -z "$seed" ] && seed=$(cat /sys/class/net/eth0/address 2>/dev/null | tr -d ':')
        [ -z "$seed" ] && seed=$(cat /proc/sys/kernel/random/boot_id 2>/dev/null | tr -d '-')
        [ -z "$seed" ] && seed="$(date +%s%N 2>/dev/null || date +%s)$$"
        local mac_raw
        mac_raw=$(printf '%s' "$seed" | md5sum | cut -c1-12)
        local b1; b1=$(printf "%02x" $(( (0x${mac_raw:0:2} & 0xFE) | 0x02 )))
        saved_mac="${b1}:${mac_raw:2:2}:${mac_raw:4:2}:${mac_raw:6:2}:${mac_raw:8:2}:${mac_raw:10:2}"
        uci set easymesh.global.agent_br_mac="$saved_mac"
        uci commit easymesh
        log "agent_prepare: generated unique MAC=%s" "$saved_mac"
    else
        log "agent_prepare: using saved MAC=%s" "$saved_mac"
    fi

    local dev_name dev_sec=""
    while IFS= read -r line; do
        case "$line" in
            network.*=device) dev_sec="${line%%=*}" ;;
        esac
        if [ -n "$dev_sec" ]; then
            local dname; dname=$(uci -q get "${dev_sec}.name" 2>/dev/null)
            if [ "$dname" = "br-lan" ]; then
                uci set "${dev_sec}.macaddr=${saved_mac}"
                uci commit network
                log "agent_prepare: MAC persisted via %s.macaddr" "$dev_sec"
                break
            fi
            dev_sec=""
        fi
    done < <(uci show network 2>/dev/null)

    local old_mac; old_mac=$(cat /sys/class/net/br-lan/address 2>/dev/null)
    if ip link set dev br-lan address "$saved_mac" 2>/dev/null; then
        log "agent_prepare: br-lan MAC %s -> %s" "$old_mac" "$saved_mac"
    else
        loge "agent_prepare: ip link set failed, trying via bridge"
        ip link set dev br-lan down 2>/dev/null
        ip link set dev br-lan address "$saved_mac" 2>/dev/null
        ip link set dev br-lan up 2>/dev/null
    fi

    log "agent_prepare: done. IP=%s MAC=%s" \
        "$(ip -4 addr show dev br-lan 2>/dev/null | awk '/inet /{print $2}' | head -1)" \
        "$(cat /sys/class/net/br-lan/address 2>/dev/null)"
}

# ── 配置 UCI wireless + dummy network.mesh section ───────────────────────────
apply_mesh_config() {
    local role backhaul mesh_id mesh_key mesh_band ssid key
    local ieee80211r mobility_domain ft_over_ds reassociation_deadline
    local ieee80211k ieee80211v routing_algo mesh_fwding

    config_load "$UCI_APP"
    config_get role                   global role                  'master'
    config_get backhaul               global backhaul              'wired'
    config_get mesh_id                global mesh_id               'OpenWrt-Mesh'
    config_get mesh_key               global mesh_key              ''
    config_get mesh_band              global mesh_band             '5g'
    config_get ssid                   global ssid                  'OpenWrt'
    config_get key                    global key                   ''
    config_get ieee80211r             global ieee80211r            '1'
    config_get mobility_domain        global mobility_domain       'aabb'
    config_get ft_over_ds             global ft_over_ds            '1'
    config_get reassociation_deadline global reassociation_deadline '1000'
    config_get ieee80211k             global ieee80211k            '1'
    config_get ieee80211v             global ieee80211v            '1'
    config_get routing_algo           global routing_algo          'BATMAN_IV'
    config_get mesh_fwding            global mesh_fwding           '0'

    cleanup_wireless

    local radio_mesh
    get_radio radio_mesh "$mesh_band"
    [ -z "$radio_mesh" ] && radio_mesh="radio0"
    log "apply_mesh_config: role=%s backhaul=%s radio=%s" "$role" "$backhaul" "$radio_mesh"

    # ── network.mesh：dummy section，只为给 wireless.mesh0.network 提供绑定 ──
    uci set network.mesh=interface
    uci set network.mesh.proto='none'
    uci commit network

    # ── 无线 mesh 回程（802.11s），backhaul=wired 时禁用 ─────────────────────
    uci set wireless.mesh0=wifi-iface
    uci set wireless.mesh0.device="$radio_mesh"
    uci set wireless.mesh0.mode='mesh'
    uci set wireless.mesh0.mesh_id="$mesh_id"
    uci set wireless.mesh0.mesh_fwding="$mesh_fwding"
    uci set wireless.mesh0.mesh_rssi_threshold='-80'
    uci set wireless.mesh0.network='mesh'
    uci set wireless.mesh0.easymesh_managed='1'
    if [ -n "$mesh_key" ]; then
        uci set wireless.mesh0.encryption='sae'
        uci set wireless.mesh0.key="$mesh_key"
    else
        uci set wireless.mesh0.encryption='none'
    fi
    uci set wireless.mesh0.disabled=$([ "$backhaul" = "wired" ] && echo '1' || echo '0')

    if [ "$role" = "master" ]; then
        uci -q del dhcp.lan.ignore
        log "Role: master — DHCP enabled"
    elif [ "$backhaul" = "wired" ]; then
        uci -q del dhcp.lan.ignore
        log "Role: agent (wired) — DHCP kept enabled"
    else
        uci set dhcp.lan.ignore='1'
        log "Role: agent (wireless) — DHCP disabled"
    fi
    uci commit dhcp

    # ── 客户端 AP ─────────────────────────────────────────────────────────────
    uci set wireless.easymesh_ap=wifi-iface
    uci set wireless.easymesh_ap.device="$radio_mesh"
    uci set wireless.easymesh_ap.mode='ap'
    uci set wireless.easymesh_ap.ssid="$ssid"
    uci set wireless.easymesh_ap.network='lan'
    uci set wireless.easymesh_ap.easymesh_managed='1'
    uci set wireless.easymesh_ap.disabled='0'
    if [ -n "$key" ]; then
        uci set wireless.easymesh_ap.encryption='psk2+ccmp'
        uci set wireless.easymesh_ap.key="$key"
    else
        uci set wireless.easymesh_ap.encryption='none'
    fi
    [ "$ieee80211r" = "1" ] && {
        uci set wireless.easymesh_ap.ieee80211r='1'
        uci set wireless.easymesh_ap.mobility_domain="$mobility_domain"
        uci set wireless.easymesh_ap.ft_over_ds="$ft_over_ds"
        uci set wireless.easymesh_ap.ft_psk_generate_local='1'
        uci set wireless.easymesh_ap.reassociation_deadline="$reassociation_deadline"
    }
    [ "$ieee80211k" = "1" ] && {
        uci set wireless.easymesh_ap.ieee80211k='1'
        uci set wireless.easymesh_ap.rrm_neighbor_report='1'
        uci set wireless.easymesh_ap.rrm_beacon_report='1'
    }
    [ "$ieee80211v" = "1" ] && {
        uci set wireless.easymesh_ap.ieee80211v='1'
        uci set wireless.easymesh_ap.bss_transition='1'
    }
    uci commit wireless
    log "apply_mesh_config done"
}

# ── procd ─────────────────────────────────────────────────────────────────────
start_service() {
    local enabled role backhaul
    config_load "$UCI_APP"
    config_get_bool enabled global enabled 0
    [ "$enabled" -eq 0 ] && return 0

    config_get role     global role     'master'
    config_get backhaul global backhaul 'wired'
    log "EasyMesh starting (role=%s)" "$role"

    apply_mesh_config

    if [ "$role" = "agent" ] && [ "$backhaul" = "wired" ]; then
        _agent_prepare
    fi

    ( sleep 2; wifi reload ) &

    log "Starting EasyMesh daemon..."
    procd_open_instance "easymesh-${role}"
    procd_set_param command /usr/sbin/easymesh-${role}
    procd_set_param respawn 10 3 0
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance
}

stop_service() {
    log "EasyMesh stopping..."
    for pid in $(pgrep -f "easymesh-(master|agent)" 2>/dev/null); do
        [ "$pid" != "$$" ] && kill "$pid" 2>/dev/null || true
    done

    cleanup_wireless

    uci -q del dhcp.lan.ignore
    uci commit dhcp

    uci -q del network.mesh
    uci commit network
    ( sleep 1; wifi reload ) &
    rm -rf /tmp/easymesh* /tmp/em_* 2>/dev/null
}

reload_service() {
    stop
    sleep 2
    start
}

service_triggers() {
    procd_add_reload_trigger "$UCI_APP"
}
