# luci-app-easymesh/root/etc/init.d/easymesh
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=95
STOP=10

. /lib/functions/network.sh

EM_RUN_DIR="/tmp/easymesh"
LOGFILE="/tmp/easymesh.log"
EM_BACKUP_DIR="/etc/easymesh/backup"
wireless_backup="${EM_BACKUP_DIR}/wireless_backup.conf"
[ -z "$lang" ] && lang=$(uci_get luci main lang en)

translate() {
    local tpl="$1"; shift
    local translated
    translated=$(lua -e "require \"luci.i18n\".setlanguage(\"${lang}\"); print(require \"luci.i18n\".translate([==[${tpl}]==]))" 2>/dev/null)
    [ -z "$translated" ] && translated="$tpl"
    [ $# -gt 0 ] && printf "$translated" "$@" || printf '%s' "$translated"
}
_emlog() {
    local level="$1" tag="$2"; shift 2
    local ts; ts=$(date "+%Y-%m-%d %H:%M:%S")
    printf "%s [%-5s] [%-8s] %s\n" "$ts" "$level" "$tag" "$*" >> "$LOGFILE"
}
log()  { _emlog "INFO"  "initd" "$(translate "$@")"; }
loge() { _emlog "ERROR" "initd" "$(translate "$@")"; }

_restore_iface() {
    local iface="$1" managed
    config_get managed "$iface" easymesh_managed
    [ "$managed" != "1" ] && return
    case "$iface" in mesh0|_5g_|_2g_|_fb_|_auto_) uci -q del wireless."$iface" ;; esac
}

cleanup_wireless() {
    config_load wireless
    config_foreach _restore_iface wifi-iface
    if [ -f "$wireless_backup" ]; then
        log "Restoring wireless config from backup"
        uci import wireless < "$wireless_backup" 2>/dev/null && rm -f "$wireless_backup"
    fi
    uci commit wireless
}

_ap_iface_result=""
_ap_target_radio=""
_find_ap_cb() {
    local n="$1" dev mode
    config_get dev  "$n" device
    config_get mode "$n" mode
    [ "$dev" = "$_ap_target_radio" ] && [ "$mode" = "ap" ] && [ "$n" != "mesh0" ] &&
        [ -z "$_ap_iface_result" ] && _ap_iface_result="$n"
}
_get_ap_iface() {
    _ap_iface_result=""; _ap_target_radio="$1"
    config_load wireless
    config_foreach _find_ap_cb wifi-iface
    printf '%s' "$_ap_iface_result"
}

_agent_setup() {
    [ "$1" != "agent" ] && return
    local ip="$2" prefix="${2%.*}" original_ip="$2"

    [ -f "${EM_BACKUP_DIR}"/easymesh_agent ] && \
        original_ip=$(sed -n 's/^ipaddr=//p' "${EM_BACKUP_DIR}"/easymesh_agent)

    [ "${ip##*.}" != "1" ] && [ "$original_ip" != "$ip" ] && return

    local mac_raw=""
    for _iface in br-lan eth0 bat0; do
        [ -f "/sys/class/net/${_iface}/address" ] && {
            mac_raw=$(cat "/sys/class/net/${_iface}/address" | tr -d ':')
            break
        }
    done
    mac_raw="${mac_raw:-000000aabbcc}"

    local mac_tail; mac_tail=$(printf '%s' "$mac_raw" | cut -c7-12)
    local offset; offset=$(printf '%d' "0x${mac_tail:-0}" 2>/dev/null || echo 0)
    offset=$((offset % 98))
    local candidate="${prefix}.$((100 + offset))"
    local new_ip="$candidate"

    if ping -c 1 -W 1 -q "$candidate" >/dev/null 2>&1; then
        new_ip=""
        local i
        for i in $(seq 100 197); do
            local c="${prefix}.${i}"
            ping -c 1 -W 1 -q "$c" >/dev/null 2>&1 || { new_ip="$c"; break; }
        done
        [ -z "$new_ip" ] && {
            loge "agent_setup: no free IP in .100-.197, keeping .1"
            return 0
        }
    fi

    [ "$new_ip" != "$ip" ] && {
        uci -q set network.lan.ipaddr="$new_ip"
        uci -q set system.@system[0].hostname="OpenWrt-Agent-${new_ip##*.}"
        uci -q commit network && uci -q commit system
        /etc/init.d/network reload
        local hostname; hostname=$(uci_get system @system[0] hostname openwrt)
        [ -s "${EM_BACKUP_DIR}"/easymesh_agent ] || \
            printf 'ipaddr=%s\nhostname=%s\n' "$ip" "$hostname" > "${EM_BACKUP_DIR}"/easymesh_agent
    }
}

_identify_radio_cb() {
    local r="$1" band
    config_get band "$r" band
    case "$band" in
        2g*|2.4*|bgn|bg|b)
            [ -z "$radio_2g" ] && radio_2g="$r" ;;
        5g*|5GHz*|an|a)
            [ -z "$radio_5g" ] && radio_5g="$r" || radio_5g2="$r" ;;
    esac
}
_identify_radios() {
    radio_2g=""; radio_5g=""; radio_5g2=""
    config_load wireless
    config_foreach _identify_radio_cb wifi-device
}

_setup_iface_result=""
_setup_ap() {
    local desired_name="$1" radio="$2" ssid="$3" key="$4" ft="$5" k="$6" v="$7" \
          domain="$8" ft_ds="$9" deadline="${10}"
    local existing; existing=$(_get_ap_iface "$radio")
    local final_name="$desired_name"

    [ -n "$existing" ] && {
        [ ! -f "$wireless_backup" ] && uci export wireless > "$wireless_backup" 2>/dev/null
        final_name="$existing"
    }

    uci -q set wireless.${final_name}.easymesh_managed=1
    uci -q set wireless.${final_name}.device="$radio"
    uci -q set wireless.${final_name}.network='lan'
    uci -q set wireless.${final_name}.mode='ap'

    [ -n "$ssid" ] && uci -q set wireless.${final_name}.ssid="$ssid"

    if [ -n "$key" ]; then
        uci -q set wireless.${final_name}.encryption='psk2+ccmp'
        uci -q set wireless.${final_name}.key="$key"
    else
        uci -q set wireless.${final_name}.encryption='none'
        uci -q del wireless.${final_name}.key
    fi

    [ "$ft" = "1" ] && {
        uci -q set wireless.${final_name}.ieee80211r=1
        uci -q set wireless.${final_name}.mobility_domain="$domain"
        uci -q set wireless.${final_name}.ft_over_ds="$ft_ds"
        uci -q set wireless.${final_name}.ft_psk_generate_local=1
        uci -q set wireless.${final_name}.reassociation_deadline="$deadline"
    }
    [ "$k" = "1" ] && {
        uci -q set wireless.${final_name}.ieee80211k=1
        uci -q set wireless.${final_name}.rrm_neighbor_report=1
        uci -q set wireless.${final_name}.rrm_beacon_report=1
    }
    [ "$v" = "1" ] && {
        uci -q set wireless.${final_name}.ieee80211v=1
        uci -q set wireless.${final_name}.bss_transition=1
    }

    _setup_iface_result="$final_name"
}

start_service() {
    config_load easymesh
    config_get_bool enabled global enabled 0
    [ "$enabled" -eq 0 ] && return 0

    config_get role        global  role                   master
    config_get backhaul    global  backhaul               wired
    config_get mesh_id     global  mesh_id                'OpenWrt-Mesh'
    config_get mesh_key    global  mesh_key               ''
    config_get mesh_band   global  mesh_band              5g
    config_get ssid        global  ssid                   'OpenWrt'
    config_get key         global  key                     ''
    config_get band_merge  global  band_merge              1
    config_get ssid_5g     global  ssid_5g                 ''
    config_get ssid_2g     global  ssid_2g                 ''
    config_get key_5g      global  key_5g                  ''
    config_get key_2g      global  key_2g                  ''
    config_get wan_mode    global  wan_mode                router
    config_get wan_proto   global  wan_proto               dhcp
    config_get wan_user    global  wan_pppoe_user          ''
    config_get wan_pass    global  wan_pppoe_pass          ''
    config_get ft          global  ieee80211r              1
    config_get domain      global  mobility_domain         aabb
    config_get ft_ds       global  ft_over_ds              1
    config_get deadline    global  reassociation_deadline  1000
    config_get k           global  ieee80211k              1
    config_get v           global  ieee80211v              1

    config_get backhaul_port global backhaul_port ''

    mkdir -p "$EM_BACKUP_DIR"
    log "EasyMesh starting (role=%s)" "$role"

    _agent_setup "$role" "$(uci_get network lan ipaddr)"
    cleanup_wireless
    _identify_radios

    local mesh_radio="$radio_5g"
    [ "$mesh_band" = "2g" ] && mesh_radio="$radio_2g"
    [ -z "$mesh_radio" ] && mesh_radio=$(uci -q show wireless 2>/dev/null | \
        awk -F= '/=wifi-device$/{sub(/wireless\./,""); print $1}' | head -1)

    uci -q set network.mesh=interface
    uci -q set network.mesh.proto=none

    if [ "$backhaul" != "wired" ] ; then
        uci -q set wireless.mesh0=wifi-iface
        uci -q set wireless.mesh0.device="$mesh_radio"
        uci -q set wireless.mesh0.mode=mesh
        uci -q set wireless.mesh0.mesh_id="$mesh_id"
        uci -q set wireless.mesh0.mesh_fwding=0
        uci -q set wireless.mesh0.mesh_rssi_threshold=-80
        uci -q set wireless.mesh0.network=mesh
        uci -q set wireless.mesh0.easymesh_managed=1
        if [ -n "$mesh_key" ]; then
            uci -q set wireless.mesh0.encryption=sae
            uci -q set wireless.mesh0.key="$mesh_key"
        else
            uci -q set wireless.mesh0.encryption=none
        fi
    fi

    if [ "$role" = "master" ]; then
        uci -q del dhcp.lan.ignore
        if [ "$band_merge" = "1" ]; then
            ssid_5g="$ssid"; key_5g="$key"; ssid_2g="$ssid"; key_2g="$key"
        else
            ssid_5g="${ssid_5g:-$ssid}"; key_5g="${key_5g:-$key}"
            ssid_2g="${ssid_2g:-$ssid}"; key_2g="${key_2g:-$key}"
        fi
    else
        uci -q set dhcp.lan.ignore=1
    fi

    local idx=0
    if [ "$role" = "master" ]; then
        [ -n "$radio_5g" ] && {
            _setup_ap "_5g_" "$radio_5g" "$ssid_5g" "$key_5g" "$ft" "$k" "$v" "$domain" "$ft_ds" "$deadline"
            uci -q del wireless.${_setup_iface_result}.disabled
            idx=$(( idx + 1 ))
        }
        [ -n "$radio_2g" ] && {
            _setup_ap "_2g_" "$radio_2g" "$ssid_2g" "$key_2g" "$ft" "$k" "$v" "$domain" "$ft_ds" "$deadline"
            uci -q del wireless.${_setup_iface_result}.disabled
            idx=$(( idx + 1 ))
        }
        [ "$idx" -eq 0 ] && {
            _setup_ap "_fb_" "${mesh_radio:-radio0}" "$ssid_5g" "$key_5g" "$ft" "$k" "$v" "$domain" "$ft_ds" "$deadline"
            uci -q del wireless.${_setup_iface_result}.disabled
        }
    else
        [ -n "$radio_5g" ] && {
            _setup_ap "_5g_" "$radio_5g" "" "" "0" "0" "0" "" "" ""
            uci -q set wireless.${_setup_iface_result}.disabled=1
            idx=$(( idx + 1 ))
        }
        [ -n "$radio_2g" ] && {
            _setup_ap "_2g_" "$radio_2g" "" "" "0" "0" "0" "" "" ""
            uci -q set wireless.${_setup_iface_result}.disabled=1
            idx=$(( idx + 1 ))
        }
        [ "$idx" -eq 0 ] && {
            _setup_ap "_fb_" "${mesh_radio:-radio0}" "" "" "0" "0" "0" "" "" ""
            uci -q set wireless.${_setup_iface_result}.disabled=1
        }
    fi
    [ -n "$radio_5g"   ] && uci -q del wireless.${radio_5g}.disabled
    [ -n "$radio_2g"   ] && uci -q del wireless.${radio_2g}.disabled
    [ -n "$mesh_radio" ] && uci -q del wireless.${mesh_radio}.disabled

    uci -q commit dhcp
    uci -q commit wireless
    uci -q commit network

    if [ "$role" = "master" ]; then
        if [ "$wan_mode" = "switch" ]; then
            uci -q set network.wan.disabled=1
            uci -q set network.wan6.disabled=1
        else
            uci -q del network.wan.disabled
            uci -q del network.wan6.disabled
            [ "$wan_proto" = "pppoe" ] && [ -n "$wan_user" ] && {
                uci -q set network.wan.proto=pppoe
                uci -q set network.wan.username="$wan_user"
                uci -q set network.wan.password="$wan_pass"
            }
        fi
        uci -q commit network
    fi

    nft insert rule inet fw4 input tcp dport 4304 accept comment EASYMESH_PORTS 2>/dev/null || true
    nft insert rule inet fw4 input udp dport 4305 accept comment EASYMESH_PORTS 2>/dev/null || true

    wifi reload
    sleep 2

    procd_open_instance "easymesh"
    procd_set_param command /usr/sbin/easymesh-"$role"
    procd_set_param env PROCD=1
    procd_set_param respawn
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param pidfile "/var/run/easymesh-${role}.pid"
    procd_close_instance
}

stop_service() {
    log "EasyMesh stopping..."

    for pid in $(pgrep -f "easymesh-(master|agent)" 2>/dev/null); do
        [ "$pid" != "$$" ] && kill "$pid" 2>/dev/null || true
    done

    [ -f ${EM_RUN_DIR}/pair_hostapd.pid ] && kill "$(cat ${EM_RUN_DIR}/pair_hostapd.pid)" 2>/dev/null
    [ -f ${EM_RUN_DIR}/pair_dnsmasq.pid ] && kill "$(cat ${EM_RUN_DIR}/pair_dnsmasq.pid)" 2>/dev/null
    iw dev easymesh-pair0 del 2>/dev/null || true
    nft -a list chain inet fw4 input 2>/dev/null | \
        awk '/4304|4305/ {print $NF}' | \
        while read handle; do
            nft delete rule inet fw4 input handle "$handle"
        done

    cleanup_wireless
    uci -q del dhcp.lan.ignore
    uci -q commit dhcp
    uci -q del network.mesh
    uci -q commit network

    ( sleep 1; wifi reload ) &
    rm -rf "$EM_RUN_DIR" /tmp/easymesh.log /tmp/easymesh.log.lock 2>/dev/null
}

reload_service() {
    stop
    sleep 2
    start
}

service_triggers() {
    procd_add_reload_trigger easymesh
}
