# luci-app-miaplus/root/etc/init.d/miaplus
#!/bin/sh /etc/rc.common

START=30
USE_PROCD=1
ipv6enable=$(uci -q get miaplus.config.ipv6enable || echo 0)
cmd="iptables"
[ "$ipv6enable" = "1" ] && cmd="ip6tables"

log_msg() {
	logger -p daemon.info -t "miaplus" -- "$*"
	echo "$*"
}

add_rule() {
	local var="$1"
	config_get_bool client "$var" enable '0'
	config_get_bool enable 'config' enable '0'
	[ "$client" -eq "0" -o "$enable" -eq '0' ] && return
	config_get_bool strict 'config' strict '0'

	$cmd -t filter -N MIAPLUS 2>/dev/null
	$cmd -C INPUT -p udp --dport 53 -j MIAPLUS >/dev/null 2>&1 || \
		$cmd -I INPUT -p udp --dport 53 -m comment --comment "miaplus_control" -j MIAPLUS
	$cmd -C INPUT -p tcp --dport 53 -j MIAPLUS >/dev/null 2>&1 || \
		$cmd -I INPUT -p tcp --dport 53 -m comment --comment "miaplus_control" -j MIAPLUS

	[ "$strict" -eq '1' ] && ! $cmd -C FORWARD -j MIAPLUS >/dev/null 2>&1 && \
		$cmd -t filter -I FORWARD -m comment --comment "miaplus_control_forward" -j MIAPLUS

	$cmd -t nat -C PREROUTING -p udp --dport 53 -j REDIRECT --to-ports 53 >/dev/null 2>&1 || \
		$cmd -t nat -A PREROUTING -p udp --dport 53 -j REDIRECT --to-ports 53 -m comment --comment "miaplus_control"
	$cmd -t nat -C PREROUTING -p tcp --dport 53 -j REDIRECT --to-ports 53 >/dev/null 2>&1 || \
		$cmd -t nat -A PREROUTING -p tcp --dport 53 -j REDIRECT --to-ports 53 -m comment --comment "miaplus_control"

	add_rule_core() {
		local type=$1 identifier=$2 match="" time_args=""

		case "$type" in
			ip|ip6) match="-s $identifier" ;;
			mac)    match="-m mac --mac-source $identifier" ;;
		esac

		[ -n "$day" ] && time_args="$time_args --monthdays $day"
		[ -n "$wds" ] && time_args="$time_args --weekdays ${wds// /,}"

		if [ "$tm" -eq "1" -a -n "$on_d" -a -n "$off_d" -a -n "$on_t" -a -n "$off_t" ]; then
			time_args="$time_args --datestart ${on_d}T${on_t} --datestop ${off_d}T${off_t}"
		else
			[ -n "$on_d" -a -n "$off_d" ] && time_args="$time_args --datestart $on_d --datestop $off_d"
			[ -n "$on_t" -a -n "$off_t" ] && time_args="$time_args --timestart $on_t --timestop $off_t"
		fi

		$cmd -t filter -A MIAPLUS $match \
			-m time --kerneltz $time_args \
			-j DROP \
			-m comment --comment "miaplus_time_${identifier}"

		$cmd -t filter -A MIAPLUS $match \
			-j ACCEPT \
			-m comment --comment "miaplus_allow_$identifier"
	}

	config_get ip "$var" ipaddr ''
	config_get mac "$var" macaddr ''
	config_get ip6 "$var" ip6addr ''
	config_get wds "$var" weekdays ''
	config_get day "$var" monthdays ''
	config_get on_d "$var" start_date ''
	config_get off_d "$var" stop_date ''
	config_get on_t "$var" start_time ''
	config_get off_t "$var" stop_time ''
	config_get_bool tm "$var" time_mode '0'

	if [ -n "$mac" ]; then
		add_rule_core "mac" "$mac"
	elif [ -n "$ip" ]; then
		add_rule_core "ip" "$ip"
	elif [ -n "$ip6" -a "$ipv6enable" = "1" ]; then
		add_rule_core "ip6" "$ip6"
	fi
}

start_service() {
	stop_service
	$cmd -t filter -N MIAPLUS 2>/dev/null
	$cmd -t filter -F MIAPLUS 2>/dev/null
	config_load miaplus
	config_foreach add_rule macbind
}

stop_service() {
	for chain in INPUT FORWARD; do
		$cmd -L $chain -n --line-numbers 2>/dev/null | \
			awk '/miaplus_/ {print $1}' | sort -rn | \
			while read num; do $cmd -D $chain $num 2>/dev/null; done
	done

	$cmd -t nat -L PREROUTING -n --line-numbers 2>/dev/null | \
		awk '/miaplus_control/ {print $1}' | sort -rn | \
		while read num; do $cmd -t nat -D PREROUTING $num 2>/dev/null; done

	$cmd -t filter -F MIAPLUS 2>/dev/null
	$cmd -t filter -X MIAPLUS 2>/dev/null
}

service_triggers() {
	procd_add_reload_trigger miaplus
}
