# luci-app-wizard/root/etc/uci-defaults/luci-app-wizard
#!/bin/sh

echo "config wizard 'default'" > /etc/config/wizard
chmod +x /etc/init.d/wizard
chmod +x /etc/hotplug.d/button/99-wizard-restore

rm -f /tmp/luci-*
exit 0
