# luci-app-wizard/root/etc/hotplug.d/button/99-wizard-restore
#!/bin/sh

[ "$BUTTON" = "reset" ] || [ "$BUTTON" = "wps" ] || exit 0
[ "$ACTION" = "released" ] || exit 0

SEEN="${SEEN:-0}"
LOCK="/var/run/wizard-restore.lock"
. /lib/config/uci.sh
led_blink() {
    local led=$(ls /sys/class/leds/*/trigger 2>/dev/null | head -1)
    led=${led%/trigger}
    [ -z "$led" ] && return
    echo timer > $led/trigger
    echo $1 > $led/delay_on
    echo $2 > $led/delay_off
    sleep $3
    echo default-on > $led/trigger 2>/dev/null
}

check_gateway() {
    local gw=$(uci_get network lan gateway "")
    [ -z "$gw" ] && gw=$(uci_get network wan gateway "")
    [ -z "$gw" ] && return 1

    ping -c 1 -W 2 "$gw" >/dev/null 2>&1
    return $?
}

[ -e "$LOCK" ] && exit 0
touch "$LOCK"
trap "rm -f '$LOCK'" EXIT

if [ "$SEEN" -ge 3 ] && [ "$SEEN" -lt 6 ]; then
    if check_gateway; then
        led_blink 200 200 0.5
    else
        led_blink 300 1000 3
    fi
    exit 0
fi

if [ "$SEEN" -ge 6 ]; then
    if check_gateway; then
        led_blink 200 200 0.5
        exit 0
    fi

    led_blink 100 100 5

    for file in dhcp network firewall wireless; do
        local bak_file="/etc/config/bk_$file"
        if [ -f "$bak_file" ]; then
            uci -q import "$cfg" < "$bak_file"
            uci -q commit "$cfg"
        fi
    done

    uci -q delete wizard.default.wan_proto 2>/dev/null
    uci -q commit wizard

    /etc/init.d/network reload
    sleep 2
    /etc/init.d/firewall reload
    /etc/init.d/dnsmasq reload
    wifi reload
fi

exit 0
