# luci-app-deluge/root/etc/uci-defaults/luci-app-deluge
#!/bin/sh
chmod +x /etc/init.d/deluge
echo "config deluge 'main'" > /etc/config/deluge

rm -rf /tmp/luci-*
exit 0
