# luci-app-easymesh/root/etc/hotplug.d/button/25-easymesh-wps
#!/bin/sh
[ "$BUTTON" = "wps" ] || exit 0
[ "$ACTION" = "released" ] || exit 0
SEEN_TIME=${SEEN:-0}
[ "$SEEN_TIME" -ge 2 ] && [ "$SEEN_TIME" -le 10 ] || exit 0

. /lib/config/uci.sh
LOCK="/var/run/easymesh-restore.lock"
[ -e "$LOCK" ] && exit 0
touch "$LOCK"
trap "rm -f '$LOCK'" EXIT

role=$(uci_get easymesh global role master)
enabled=$(uci_get easymesh global enabled 0)
if [ "$enabled" = "1" ] && [ "$role" = "agent" ] && [ -f /etc/easymesh_agent ]; then
    original_ip=$(sed -n 's/^ipaddr=//p' /etc/easymesh_agent)
    original_hostname=$(sed -n 's/^hostname=//p' /etc/easymesh_agent)
    /etc/init.d/easymesh stop
    if [ -n "$original_hostname" ]; then
        uci_set system @system[0] hostname "$original_hostname"
        uci_commit system
    fi
    if [ -n "$original_ip" ]; then
        uci_set network lan ipaddr "$original_ip"
        uci_commit network
        /etc/init.d/network reload &
    fi

    uci_set easymesh global enabled 0
    uci_remove easymesh global role
    uci_commit easymesh
    ( sleep 1; /etc/init.d/easymesh restart ) &
    rm -f /etc/easymesh_agent
elif [ "$enabled" != "1" ] && [ "$role" != "agent" ] && [ ! -f /etc/easymesh_agent ] ; then
    uci_set easymesh global enabled 1
    uci_set easymesh global role agent
    uci_commit easymesh
    ( sleep 1; /etc/init.d/easymesh restart ) &
else
    exit 0
fi

(
    RED_LED=""; GRN_LED=""
    for led in $(ls /sys/class/leds 2>/dev/null); do
        [ -z "$RED_LED" ] && echo "$led" | grep -Eiq "red|fault|err" && [ -w "/sys/class/leds/$led/brightness" ] && RED_LED="/sys/class/leds/$led/brightness"
        [ -z "$GRN_LED" ] && echo "$led" | grep -Eiq "green|blue|white|status|system" && [ -w "/sys/class/leds/$led/brightness" ] && GRN_LED="/sys/class/leds/$led/brightness"
    done

    set_ms_blink() {
        local path="$1" on="$2" off="$3"
        [ -z "$path" ] && return
        local base="${path%/*}"
        if [ "$on" -gt 0 ]; then
            echo "timer" > "$base/trigger" 2>/dev/null
            echo "$on" > "$base/delay_on" 2>/dev/null
            echo "$off" > "$base/delay_off" 2>/dev/null
        else
            echo "none" > "$base/trigger" 2>/dev/null
            echo 0 > "$path" 2>/dev/null
        fi
    }

    if [ -n "$RED_LED" ] && [ -n "$GRN_LED" ]; then
        for i in 1 2 3 4 5 6; do
            set_ms_blink "$RED_LED" 300 300;
            set_ms_blink "$GRN_LED" 0 0; sleep 0.5
            set_ms_blink "$RED_LED" 0 0;
            set_ms_blink "$GRN_LED" 300 300; sleep 0.5
        done
        set_ms_blink "$RED_LED" 0 0
        echo "default-on" > "${GRN_LED%/*}/trigger" 2>/dev/null || echo 1 > "$GRN_LED" 2>/dev/null
    elif [ -n "$RED_LED" ] || [ -n "$GRN_LED" ]; then
        TARGET=${RED_LED:-$GRN_LED}
        set_ms_blink "$TARGET" 300 300
        sleep 3
        echo "default-on" > "${TARGET%/*}/trigger" 2>/dev/null || echo 1 > "$TARGET" 2>/dev/null
    fi
) > /dev/null 2>&1 &

exit 0
