# luci-app-easymesh/root/etc/hotplug.d/button/99-easymesh-wps
#!/bin/sh
# /etc/hotplug.d/button/99-easymesh-wps
#
# 按下 WPS/Reset 键 2~10 秒松开 → 以 agent 身份加入 Mesh
# 适用于：新子节点首次接入，无需打开 LuCI 界面
# 条件：当前 easymesh 未启用（enabled=0）才触发，防止误操作

# 只响应 WPS / Reset / BTN_0 / BTN_1 按钮
case "$BUTTON" in
    wps|reset|BTN_0|BTN_1) ;;
    *) exit 0 ;;
esac

# 松开事件 + 按住时间 2~10 秒
[ "$ACTION" = "released" ] || exit 0
[ -n "$SEEN" ] || exit 0
[ "$SEEN" -ge 2 ] && [ "$SEEN" -le 10 ] || exit 0

# 仅在未启用时触发（防止已配置节点被重置）
enabled=$(uci -q get easymesh.global.enabled)
[ "$enabled" = "1" ] && exit 0

# 配置为子节点并启用
uci set easymesh.global.role='agent'
uci set easymesh.global.enabled='1'
uci commit easymesh

logger -t EasyMesh "WPS button: enabling as agent, restarting service"

# 重启服务（后台执行，不阻塞 hotplug）
( sleep 1; /etc/init.d/easymesh restart ) &

# LED 反馈：快速闪烁 3 次（兼容各厂商 LED 命名）
LED_PATH=""
for led_pattern in wps mesh system status power wan; do
    LED_PATH=$(ls /sys/class/leds/*${led_pattern}*/brightness 2>/dev/null | head -1)
    [ -n "$LED_PATH" ] && [ -w "$LED_PATH" ] && break
    LED_PATH=""
done

if [ -n "$LED_PATH" ]; then
    for i in 1 2 3; do
        echo 1 > "$LED_PATH" 2>/dev/null; sleep 0.2
        echo 0 > "$LED_PATH" 2>/dev/null; sleep 0.2
    done
    # 最终点亮，表示服务正在启动中
    echo 1 > "$LED_PATH" 2>/dev/null
fi
