# luci-app-easymesh/root/etc/hotplug.d/button/99-easymesh-wps
#!/bin/sh

case "$BUTTON" in
    wps|reset|BTN_0|BTN_1|rfkill|wlan) ;;
    *) exit 0 ;;
esac

[ "$ACTION" = "released" ] || exit 0
SEEN_TIME=${SEEN:-0}
[ "$SEEN_TIME" -ge 2 ] && [ "$SEEN_TIME" -le 10 ] || exit 0

enabled=$(uci -q get easymesh.global.enabled)
if [ "$enabled" = "1" -a -f /tmp/easymesh_agent_ip ]; then
    original_ip=$(cat /tmp/easymesh_agent_ip)
    uci set network.lan.ipaddr="$original_ip"
    uci commit network

    /etc/init.d/network reload
    sleep 2
else
    uci set easymesh.global.role='agent'
    uci set easymesh.global.enabled='1'
    uci commit easymesh
    ( sleep 1; /etc/init.d/easymesh restart ) &
fi

(
    RED_LED=""
    GRN_LED=""

    for led in $(ls /sys/class/leds 2>/dev/null); do
        if [ -z "$RED_LED" ] && echo "$led" | grep -Ei "red|fault|err"; then
            [ -w "/sys/class/leds/$led/brightness" ] && RED_LED="/sys/class/leds/$led/brightness"
        fi
        if [ -z "$GRN_LED" ] && echo "$led" | grep -Ei "green|blue|white|status|system"; then
            [ -w "/sys/class/leds/$led/brightness" ] && GRN_LED="/sys/class/leds/$led/brightness"
        fi
    done

    if [ -n "$RED_LED" ] && [ -n "$GRN_LED" ]; then
        RED_MAX=$(cat "${RED_LED%/*}/max_brightness" 2>/dev/null || echo 1)
        GRN_MAX=$(cat "${GRN_LED%/*}/max_brightness" 2>/dev/null || echo 1)

        for i in 1 2 3 4 5 6; do
            echo "$RED_MAX" > "$RED_LED" 2>/dev/null
            echo 0 > "$GRN_LED" 2>/dev/null
            sleep 1
            echo 0 > "$RED_LED" 2>/dev/null
            echo "$GRN_MAX" > "$GRN_LED" 2>/dev/null
            sleep 1
        done
        echo 0 > "$RED_LED" 2>/dev/null
        echo "$GRN_MAX" > "$GRN_LED" 2>/dev/null

    elif [ -n "$RED_LED" ] || [ -n "$GRN_LED" ]; then
        TARGET_LED=${RED_LED:-$GRN_LED}
        MAX=$(cat "${TARGET_LED%/*}/max_brightness" 2>/dev/null || echo 1)
        for i in 1 2 3 4 5 6; do
            echo "$MAX" > "$TARGET_LED" 2>/dev/null; sleep 1
            echo 0 > "$TARGET_LED" 2>/dev/null; sleep 1
        done
        echo "$MAX" > "$TARGET_LED" 2>/dev/null
    fi
) &

exit 0
