# luci-app-deluge/root/etc/init.d/deluge
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

log_msg() {
    logger -p daemon.info -t "deluge" "$*"
}

get_config() {
    for val in $@; do
        config_get "$val" main "$val" ""
    done
}

create_directory() {
    for dir in "$@"; do
        [ -n "$dir" -a ! -d "$dir" ] && mkdir -p "$dir" 2>/dev/null
    done
}

set_ownership() {
    [ -z "$user" ] && return 0
    id "$user" >/dev/null 2>&1 || return 1

    _paths="$log_file $log_dir $profile_dir $download_location $move_completed_path $torrentfiles_location"

    echo "$_paths" | while read -r path; do
        [ -z "$path" ] && continue
        if [ -d "$path" ]; then
            chown -R "$user" "$path"
            chmod 750 "$path"
        elif [ -f "$path" ]; then
            chown "$user" "$path"
            chmod 640 "$path"
        fi
    done
}

start_deluge() {
    config_get_bool enabled main enabled 0
    [ "$enabled" -eq 0 ] && return 1

    get_config  port user log_dir log_level profile_dir enable_logging \
                download_location move_completed_path torrentfiles_location
    create_directory "$profile_dir" "$download_location"

    local log_param log_file
    local port_param=${port:+"-p $port"}
    local user_param=${user:+"-U $user"}
    if [ "$enable_logging" -eq 1 ]; then
        create_directory "$log_dir"
        log_file="${log_dir:-$profile_dir}/deluge.log"
        touch "$log_file"
        log_param="-l $log_file ${log_level:+"-L $log_level"}"
    fi

    set_ownership
    /usr/bin/python3 /usr/libexec/update_deluge.py

    procd_open_instance 'deluged'
    procd_set_param command /usr/bin/deluged -d -c $profile_dir $log_param $user_param
    procd_set_param respawn 300 5 5
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param file /usr/bin/deluged
    procd_set_param term_timeout 10
    procd_close_instance

    procd_open_instance 'deluge-web'
    procd_set_param command /usr/bin/deluge-web -d -c $profile_dir $log_param $user_param $port_param
    procd_set_param respawn
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param file /usr/bin/deluge-web
    procd_set_param term_timeout 10
    procd_close_instance

    (
        sleep 5
        for proc in deluged deluge-web; do
            pgrep -x "$proc" >/dev/null \
                && log_msg "$proc is running" \
                || log_msg "$proc failed to start"
        done
    ) &
}

start_service() {
    config_load deluge
    config_foreach start_deluge
}

stop_service() {
    local deluge_pid=$(pgrep -x deluged)
    [ -z "$deluge_pid" ] && return 0

    kill -TERM "$deluge_pid" 2>/dev/null

    timeout=10
    while kill -0 "$deluge_pid" 2>/dev/null && [ $((timeout--)) -gt 0 ]; do
        sleep 1
    done

    kill -0 "$deluge_pid" 2>/dev/null && kill -KILL "$deluge_pid"
    log_msg 'Stop deluged service'
}

reload_service() {
    stop
    start
}

service_triggers() {
    procd_add_reload_trigger deluge
}
