# luci-app-taskplan/root/etc/uci-defaults/luci-taskplan
#!/bin/sh
mkdir -p /etc/taskplan
chmod +x /etc/init.d/taskplan \
		 /usr/bin/taskplanhandler
rm -rf /tmp/luci-*
exit 0
